# phenix-proxy.service
[Unit]
Description=PhenixRTS CORS Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/web
ExecStart=/usr/bin/python3 /opt/web/proxy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
